# bin/conf/supervisor/nfschk.conf
[program:nfschk]
command =bash /root/agent/sh/nfschk.sh
autostart = true
autorestart = true
user = root