# bin/conf/supervisor/vncc.conf
[program:vncc]
command =bash /root/vncc/vncc.sh
autostart = true
autorestart = true
user = root
