# bin/conf/supervisor/vxagent.conf
[program:vxagent]
command = bash /root/agent/vxagent.sh
autostart = true
autorestart = true
startretries=50
user = root